# server/supervisor.conf
[program:aiwiki]
command=/srv/aiwiki/run.sh
numprocs=1
autostart=true
autorestart=true
user=root
redirect_stderr=true
stdout_logfile=/var/log/supervisor/aiwiki